# src/etc/cron.d/yacd
*/8 * * * * nobody /opt/yacd/bin/yacg >/dev/null 2>&1
